# .config/hypr/systemd/hypr-battery-monitor.service
[Unit]
Description=Hyprland Battery Monitor Check
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/bin/hypr-battery-monitor.sh
Environment=DISPLAY=:0
LogLevelMax=warning
